# all/etc/systemd/system/magisystem.target
[Unit]
Description=Magisystem Services

[Install]
WantedBy=default.target
