# all/etc/systemd/system/podman-image-prune.service
[Unit]
Description=Daily Podman image prune
Before=podman-auto-update.service

[Service]
Type=oneshot
ExecStart=podman image prune -f
