# all/etc/systemd/system/wg-quick-build.timer
[Unit]
Description=Daily Wireguard Image build

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
