# balthasar/etc/systemd/system/nextcloud-cron.service
[Unit]
Description=Nextcloud background job scheduler

[Service]
Type=oneshot
ExecStart=/usr/bin/podman exec --user www-data nextcloud php cron.php
