# ayanami/etc/systemd/system/beets-import.timer
[Unit]
Description=Import music with Beets

[Timer]
OnActiveSec=0
OnUnitInactiveSec=10m
AccuracySec=1m

[Install]
WantedBy=magisystem.target
