# ayanami/etc/systemd/system/update-dyndns.timer
[Unit]
Description=Trigger DNS Update

[Timer]
OnActiveSec=0
OnUnitInactiveSec=30m

[Install]
WantedBy=network-online.target
